# funnelsort/src/profile.d/funnel.sh
alias space='du -xb * | funnel -k[i -o0 desc] | head'